# Backend/glyphmart-backend.service
[Unit]
Description=GlyphMart Flask Backend
After=network.target

[Service]
Type=exec
User=www-data
Group=www-data
WorkingDirectory=/path/to/your/GlyphMart/Backend
Environment=PATH=/path/to/your/GlyphMart/Backend/venv/bin
ExecStart=/path/to/your/GlyphMart/Backend/venv/bin/gunicorn --bind 127.0.0.1:5000 --workers 2 --worker-class sync app:app
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
